# root/etc/services.d/bind/run
#!/usr/bin/execlineb -P

/usr/sbin/named -f -c /config/named.conf -L /log/named.log
